# roles/ruby-sky2/files/home/k0kubun/.config/systemd/user/skybench.service
[Unit]
Description=Run Skybench

[Service]
Type=simple
ExecStartPre=/usr/bin/git -C /home/k0kubun/skybench pull --rebase origin master
ExecStart=/bin/bash -c 'export PATH="/home/k0kubun/.rbenv/bin:${PATH}"; exec /home/k0kubun/skybench/bin/run'

[Install]
WantedBy=default.target
