# roles/ruby-sky2/files/home/k0kubun/.config/systemd/user/skybench.timer
[Unit]
Description=Timer to run Skybench

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=skybench.service

[Install]
WantedBy=timers.target
