# roles/ruby-sky2/files/lib/systemd/system/sky2-bench.service
[Unit]
Description=Run sky2-bench

[Service]
User=k0kubun
Type=simple
ExecStartPre=/bin/bash -c 'cd /home/k0kubun/sky2-bench; git fetch origin master; git reset --hard remotes/origin/master'
ExecStart=/bin/bash -c 'export PATH="/home/k0kubun/.rbenv/bin:${PATH}"; exec /home/k0kubun/sky2-bench/bin/sky2-bench.sh'

[Install]
WantedBy=default.target
