# roles/ruby-sky2/files/lib/systemd/system/sky2-bench.timer
[Unit]
Description=Timer to run sky2-bench

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
Unit=sky2-bench.service

[Install]
WantedBy=timers.target
